# debian/cron.d
*/5 * * * * root /usr/sbin/sensor-update

